# conf/flipd.socket
[Socket]
ListenStream=/run/foobar.sk

[Install]
WantedBy=sockets.target